# ui.service
[Unit]
After=network.target

[Service]
ExecStart=/etc/systemd/system/ui.sh

[Install]
WantedBy=default.target
