# ui.service
[Unit]
After=network.target

[Service]
ExecStart=/etc/wireguard/

[Install]
WantedBy=default.target
